# deploy/aemo-negative-prices.timer
[Unit]
Description=Monitor AEMO Negative Prices for updated monthly data
Documentation=https://github.com/cutout-z/aemo-negative-prices

[Timer]
# Check daily so delayed AEMO DISPATCHPRICE archives are picked up soon after publication.
OnCalendar=*-*-* 00:35:00
RandomizedDelaySec=20m
Persistent=true
Unit=aemo-negative-prices.service

[Install]
WantedBy=timers.target
